# sh/etc/hotplug.d/net/99-meta-route
#!/bin/sh
# OpenWrt hotplug 入口：TUN 接口 add 事件触发时，调用统一的路由修正脚本。
# 真正的逻辑（删 mihomo /30 错路由 + 加正确网段）在 sh/setup-fake-ip-route.sh，
# 这一份是为了让 systemd 与 hotplug 共享同一段实现。

[ "$ACTION" = "add" ] || exit 0

# 优先用部署在 /etc/proxy/sh 的副本（OpenWrt 默认安装路径）
for _p in /etc/proxy/sh/setup-fake-ip-route.sh \
          /usr/local/lib/myrule/setup-fake-ip-route.sh; do
    if [ -x "$_p" ]; then
        exec "$_p"
    fi
done

# 兜底：找不到统一脚本时，沿用旧实现
ip route del 172.16.0.0/30 2>/dev/null
ip route replace 172.16.0.0/12 dev "$INTERFACE"
logger -t MyProxy "fake-ip fallback route -> $INTERFACE"
