# sh/etc/hotplug.d/net/99-meta-route
#!/bin/sh
# OpenWrt hotplug 入口：仅当 TUN_IFACES 列表中之一的接口 add 时，调用统一脚本。
# 其他接口（USB 网卡、WiFi roam、wg 隧道等）的 add 事件直接忽略。

# 仅对 add 事件响应
[ "$ACTION" = "add" ] || exit 0

# 加载 TUN_IFACES 等配置
# shellcheck disable=SC1091
[ -f /etc/proxy/sh/env.conf ] && . /etc/proxy/sh/env.conf
# 兜底默认值
TUN_IFACES="${TUN_IFACES:-utun Meta}"

# $INTERFACE 必须命中 TUN_IFACES，否则不动
match=0
for i in $TUN_IFACES; do
    [ "$INTERFACE" = "$i" ] && match=1 && break
done
[ "$match" = "1" ] || exit 0

# 命中则交给统一脚本处理（hotplug 触发时 TUN 与 /30 都已就绪）
exec /etc/proxy/sh/setup-fake-ip-route.sh
