# sh/etc/hotplug.d/net/99-meta-route
#!/bin/sh

[ "$ACTION" = "add" ] && { [ "$INTERFACE" = "Meta" ] || [ "$INTERFACE" = "utun" ]; } && {
    ip route del 172.16.0.0/30
    ip route add 172.16.0.0/12 dev "$INTERFACE"
    logger "Route added to $INTERFACE via hotplug"
}
