# sh/etc/systemd/system/agh.service
[Unit]
Description=AdGuardHome (MyProxy 部署)
Documentation=https://github.com/AfxMsgBox/MyRule
Wants=network-online.target
After=network-online.target proxy_core.service
StartLimitIntervalSec=3600
StartLimitBurst=5

[Service]
Type=simple
# AGH 二进制约定在 /usr/bin/AdGuardHome；如需改，写到 env.local.conf 的 AGH_BIN 即可。
ExecStart=/bin/sh -c '. /etc/proxy/sh/env.conf && exec "${AGH_BIN:-/usr/bin/AdGuardHome}" -c "$AGH_DIR/agh.yaml" -w "$AGH_DIR" --no-check-update'

Restart=on-failure
RestartSec=5

# AGH 需要 53 端口
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW

StandardOutput=journal
StandardError=journal
SyslogIdentifier=MyProxy

[Install]
WantedBy=multi-user.target
