# sh/etc/systemd/system/proxy_core.service
[Unit]
Description=MyProxy core (mihomo / clash 兼容)
Documentation=https://github.com/AfxMsgBox/MyRule
Wants=network-online.target
After=network-online.target nss-lookup.target
Before=adguardhome.service agh.service
# 与 OpenWrt procd respawn 3600 5 5 等价：1 小时窗口内崩溃 5 次后停止重启
StartLimitIntervalSec=3600
StartLimitBurst=5

[Service]
Type=simple
# 用 sh 引导以 source env.conf 中的 CORE_BIN / CORE_DIR；env.conf 包含 ${VAR:-default}
# 形式的展开，systemd 的 EnvironmentFile 解析不了。
ExecStart=/bin/sh -c '. /etc/proxy/sh/env.conf && exec "$CORE_BIN" -d "$CORE_DIR"'
# mihomo 启动后会创建 TUN 接口；这一步等接口就绪后修正路由（参见 setup-fake-ip-route.sh）。
ExecStartPost=/bin/sh -c '/etc/proxy/sh/setup-fake-ip-route.sh || true'

Restart=on-failure
RestartSec=5

# 内核要管理 TUN，需要 NET_ADMIN；按需放开
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW

# 日志走 journald（journalctl -u proxy_core -f）
StandardOutput=journal
StandardError=journal
SyslogIdentifier=MyProxy

# 防止失控：不允许内核 fork 子守护或长期占用临时目录
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
