# etc/hotplug.d/net/99-meta-route
#!/bin/sh

[ "$ACTION" = "add" ] && [ "$INTERFACE" = "Meta" ] && {
    ip route del 172.16.0.0/30
    ip route add 172.16.0.0/12 dev Meta
    logger "Route added to Meta via hotplug"
}
